# backups/docs_site/20260317_160832_pre_major_overhaul/docs/assets/systemd/11F_140/etc/systemd/system/cpe-metrics-agent.service
[Unit]
Description=CPE Metrics Agent (Serial) -> API
After=network.target
Wants=charter-api.service
After=charter-api.service
[Service]
Type=simple
User=root
WorkingDirectory=/home/da40/charter
EnvironmentFile=-/etc/default/cpe-metrics-agent
ExecStart=/home/da40/charter/tools/run_cpe_metrics_agent.sh
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
