# docs/charter_site/docs/assets/systemd/11F_140/etc/systemd/system/charter-api.service.d/10-db.conf
[Service]
Environment=DATABASE_URL=postgresql+psycopg2://rg:rg@127.0.0.1:5432/rg
Environment=PSQL_URL=postgresql://rg:rg@127.0.0.1:5432/rg
